# run.conf
#!upstart
description "cousera"

start on (local-filesystems and net-device-up)
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

exec npm start >> /var/log/myapp.log 2>&1